# infrastructure/config/upstart-minecraftd.conf
start on stopped rc RUNLEVEL=3
stop on stopping rc RUNLEVEL=[!3]

task
setuid minecraftd
exec /srv/minecraft/start.sh
